# extra/conf/bridge-socket.conf.in
description "bridge-socket: emit socket events"
emits socket

start on stopped net-loopback RESULT='ok'
stop on starting sys-down

expect readyfd
respawn
exec @startupbridgedir@/socket --readyfd
