# extra/conf/bridge-dbus.conf.in
description "bridge-dbus: emit system bus signal events"

emits dbus

start on startup

stop on runlevel [!2345]

expect readyfd
respawn

exec @startupbridgedir@/dbus --readyfd --system
