# extra/conf/bridge-udev.conf.in
description "bridge-udev: emit device events"
emits *-device-added
emits *-device-removed
emits *-device-changed
# http://www.kernel.org/pub/linux/utils/kernel/hotplug/libudev/libudev-udev-device.html
emits *-device-online
emits *-device-offline
# New in Linux 4.12+
emits *-device-bind
emits *-device-unbind

start on starting udevd
stop on stopped udevd

expect readyfd
respawn
exec @startupbridgedir@/udev --readyfd
