# extra/conf/bridge-file.conf.in
description "bridge-file: emit file events"
emits file

# the bridge does not currently handle dealing with mounts that overlay
# already-watched directories.
start on stopped fs-up RESULT='ok'
stop on starting sys-down

expect readyfd
respawn
exec @startupbridgedir@/file --readyfd
