# extra/conf-session/bridge-socket.conf.in
description "bridge-socket: emit socket events"

emits socket

start on startup
stop on session-end

expect daemon
respawn

exec @startupbridgedir@/socket --user --daemon
