# extra/conf-session/bridge-udev.conf.in
description "bridge-udev: emit device events"

# From udev itself
emits *-device-added
emits *-device-removed
emits *-device-changed
# From http://www.kernel.org/pub/linux/utils/kernel/hotplug/libudev/libudev-udev-device.html
emits *-device-online
emits *-device-offline
# New in Linux 4.12+
emits *-device-bind
emits *-device-unbind

start on startup
stop on session-end

expect readyfd
respawn

exec @startupbridgedir@/udev --user --readyfd
