# extra/conf-session/bridge-event.conf.in
description "bridge-event: emit system events in session"

emits :sys:*

start on started dbus
stop on stopped dbus

respawn

exec @startupbridgedir@/event
