# extra/conf-session/bridge-file.conf.in
description "bridge-file: emit file events"

emits file

start on started dbus
stop on session-end

expect daemon
respawn

exec @startupbridgedir@/file --daemon --user
