# extra/conf-session/bridge-dbus-session.conf.in
description "bridge-dbus-session: emit session bus signal events"
emits dbus

start on started dbus
stop on stopping dbus

expect readyfd
respawn
exec @startupbridgedir@/dbus --readyfd --session --user --bus-name session
