# extra/conf-session/dconf.conf.in
description "bridge/dconf: emit dconf events"

emits dconf

start on started dbus
stop on stopped dbus

respawn

exec @startupbridgedir@/dconf
