# extra/conf-session/re-exec.conf
description "Restart session init when system init restarts"
author "Stéphane Graber <stgraber@ubuntu.com>"
start on :sys:restarted
task
exec initctl re-exec
