# extra/conf-session/udev.conf.in
description "bridge/udev: emit device events"

# From udev itself
emits *-device-added
emits *-device-removed
emits *-device-changed
# From http://www.kernel.org/pub/linux/utils/kernel/hotplug/libudev/libudev-udev-device.html
emits *-device-online
emits *-device-offline

start on startup
stop on session-end

expect daemon
respawn

exec @startupbridgedir@/udev --user --daemon
