# util/conf/initctl-notify-dbus.conf
start on started dbus
task
exec initctl notify-dbus-address ""
