# util/conf/initctl-notify-disk.conf
start on stopped fs-up RESULT='ok'
task
pre-start exec mkdir -p /var/log/startup
exec initctl notify-disk-writeable
